# etc/udev/rules.d/50-schwinn810.rules
SUBSYSTEMS=="usb",ATTRS{idVendor}=="10c4",ATTRS{idProduct}=="ea61",SYMLINK+="schwinn810",RUN+="/usr/local/bin/schwinn810"

